# config/k3s/containerd-cleanup.timer
[Unit]
Description=Periodically clean stale k3s containerd state

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
RandomizedDelaySec=30s
Persistent=true
Unit=k3s-containerd-cleanup.service

[Install]
WantedBy=timers.target
